# scripts/com.opencli.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opencli.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>~/.nvm/versions/node/v20.20.2/bin/opencli</string>
        <string>doctor</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:~/.nvm/versions/node/v20.20.2/bin</string>
        <key>HOME</key>
        <string>~</string>
    </dict>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/Polarisor/digist/data/logs/opencli-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>~/Polarisor/digist/data/logs/opencli-daemon.log</string>
</dict>
</plist>
